# puppet/supervisor/conf.d/analytics-ui-uat.conf
; /etc/supervisor/conf.d/analytics-ui-uat.conf
;
; This file is managed by Puppet, do not edit locally !

[program:analytics-ui]
autorestart = true
autostart = true
directory = /opt/analytics-report-builder
command = /usr/local/bin/node ./src_output/server/main.js
log_stdout=true
log_stderr=true
logfile=/var/log/supervisor/ui.out
logfile_maxbytes=20MB
logfile_backups=10
startsecs=10
startretries=2
priority=999
stopasgroup=true
user=netmining
environment=NODE_ENV="UAT",LANG="en_US.UTF-8",AWS_ACCESS_KEY_ID="AKIAIQFCW7J3OPM5SQUA",AWS_SECRET_ACCESS_KEY="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"
