# puppet/supervisor/conf.d/dspot.conf
; /etc/supervisor/conf.d/dspot.conf
;
; This file is managed by Puppet, do not edit locally !

[program:dspot-api]
autorestart = true
autostart = true
command = /usr/bin/python /var/ignitionone/projects/rserver/api/prod/spot_api.py
log_stdout=true
log_stderr=true
logfile=/var/log/ignitionone/dspot-api.out
logfile_maxbytes=20MB
logfile_backups=10
startsecs=10
startretries=999
user=ignitionone

[program:dspot-dispatcher]
autorestart = true
autostart = true
command = /usr/bin/python /var/ignitionone/projects/rserver/api/prod/dispatcher.py
log_stdout=true
log_stderr=true
logfile=/var/log/ignitionone/dspot-dispatcher.out
logfile_maxbytes=20MB
logfile_backups=10
startsecs=10
startretries=999
user=ignitionone
