# puppet/supervisor/conf.d/exportq-worker.conf
; /etc/supervisor/conf.d/exportq-worker.conf
;
; This file is managed by Puppet, do not edit locally !
<%- if ( @fqdn != 'vadmin-onsite-001.us-atl.ignitionone.com' ) -%>
<%- case @area
when "production", "staging"
-%>
[program:exportq-worker]
autorestart = true
autostart = true
directory = /home/netmining
numprocs = 5
process_name = exportq_worker_%(process_num)s
environment = SUPERVISOR_PROCESS='worker%(process_num)s'
command = /usr/bin/php /var/netmining/projects/exportq/export-worker/bin/export.php --action process --conf /var/netmining/projects/exportq/export-worker/conf/netmining-activemq-stomp-production-<%= @location %>.inc.php --verbose
user=netmining
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=1
<%- else -%>
[program:exportq-preprod-worker]
autorestart = true
autostart = false
directory = /home/netmining
numprocs = 3
process_name = exportq_preprod_worker_%(process_num)s
environment = SUPERVISOR_PROCESS='worker%(process_num)s'
command = /usr/bin/php /var/netmining/projects/exportq/export-worker/bin/export.php --action process --conf /var/netmining/projects/exportq/export-worker/conf/netmining-activemq-stomp-<%= @area %>-<%= @location %>.inc.php --verbose
user=netmining
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=1
<%- end -%>
<%- end -%>

