# puppet/supervisor/conf.d/kafka-mirror-fra.conf
; /etc/supervisor/conf.d/kafka-supervisor-mirror--production.conf
;
; This file is managed by Puppet, do not edit locally !
[group:kafka-mirror-fra]
programs=kafka-tagman_audiences-fra
priority=999
[program:kafka-tagman_audiences-fra]
autorestart = true
autostart = true
directory = /usr/bin
command = kafka-mirror-maker --consumer.config /etc/kafka/consumer_fra_production.properties --producer.config /etc/kafka/producer_production.properties  --whitelist 'tagman_audiences' --num.streams 2
log_stdout=true
log_stderr=true
logfile=/var/log/kafka/kafkamirror_tagman_audiences_fra.out
logfile_maxbytes=20MB
logfile_backups=10
startsecs=10
startretries=2
priority=999
stopasgroup=true
user=netmining
environment=KAFKA_JMX_OPTS="-Dcom.sun.management.jmxremote.port=2100 -Dcom.sun.management.jmxremote.authenticate=false -Dcom.sun.management.jmxremote.ssl=false"
