# puppet/supervisor/conf.d/kafka-supervisor-rest.conf
; /etc/supervisor/conf.d/kafka-supervisor-rest.conf
;
; This file is managed by Puppet, do not edit locally !

[program:kafka-supervisor-rest]
autorestart = true
autostart = true
directory = /usr/bin
command = kafka-rest-start /etc/kafka-rest/kafka-rest.properties
log_stdout=true
log_stderr=true
logfile=/var/log/kafka/kafkarestServer.out
logfile_maxbytes=20MB
logfile_backups=10
startsecs=10
startretries=2
priority=999
stopasgroup=true
user=netmining
<%- if @location == 'us-wdc' -%>
environment=KAFKAREST_HEAP_OPTS="-Xmx2g -Xms2g",KAFKAREST_JMX_OPTS="-Dcom.sun.management.jmxremote.port=1098 -Dcom.sun.management.jmxremote.authenticate=false -Dcom.sun.management.jmxremote.ssl=false"
<%- elsif @location != 'us-atl' -%>
environment=KAFKAREST_HEAP_OPTS="-Xmx1g -Xms1g",KAFKAREST_JMX_OPTS="-Dcom.sun.management.jmxremote.port=1098 -Dcom.sun.management.jmxremote.authenticate=false -Dcom.sun.management.jmxremote.ssl=false"
<%- else -%>
environment=KAFKAREST_JMX_OPTS="-Dcom.sun.management.jmxremote.port=1098 -Dcom.sun.management.jmxremote.authenticate=false -Dcom.sun.management.jmxremote.ssl=false"
<%- end -%>
