# puppet/supervisor/conf.d/kafka-supervisor.conf
; /etc/supervisor/conf.d/kafka-supervisor.conf
;
; This file is managed by Puppet, do not edit locally !

[program:kafka-supervisor]
autorestart = true
autostart = true
directory = /usr/bin
command = kafka-server-start /etc/kafka/server.properties
log_stdout=true
log_stderr=true
logfile=/var/log/kafka/kafkaServer.out
logfile_maxbytes=20MB
logfile_backups=10
startsecs=10
startretries=2
stopasgroup=true
priority=1
user=netmining
<% if @location != 'us-atl' -%>
environment=KAFKA_HEAP_OPTS="-Xmx6g -Xms6g",KAFKA_JMX_OPTS="-Dcom.sun.management.jmxremote.port=1099 -Dcom.sun.management.jmxremote.authenticate=false -Dcom.sun.management.jmxremote.ssl=false"
<%- else -%>
environment=KAFKA_JMX_OPTS="-Dcom.sun.management.jmxremote.port=1099 -Dcom.sun.management.jmxremote.authenticate=false -Dcom.sun.management.jmxremote.ssl=false"
<%- end -%>
