# puppet/supervisor/conf.d/kafka-manager.conf
; /etc/supervisor/conf.d/kafka-manager.conf
;
; This file is managed by Puppet, do not edit locally !

[program:kafka-manager]
autorestart = true
autostart = true
directory = /usr/bin
command = kafka-manager -Dhttp.port=9002
log_stdout=false
log_stderr=true
logfile=/var/log/kafka-manager/supervisor.out
logfile_maxbytes=20MB
logfile_backups=10
startsecs=10
startretries=2
priority=999
stopasgroup=true
user=netmining
