# puppet/supervisor/conf.d/kafka-mirror-wdc.conf

; /etc/supervisor/conf.d/kafka-supervisor-mirror--production.conf
;
; This file is managed by Puppet, do not edit locally !
<%- if @location == 'zg-ash' -%>
[group:kafka-mirror-wdc]
programs=kafka-datacollection-wdc,kafka-production-wdc,kafka-video_event_json-wdc
priority=999
[program:kafka-datacollection-wdc]
autorestart = true
autostart = true
directory = /usr/bin
command = kafka-mirror-maker --consumer.config /etc/kafka/consumer_newconfig.properties --producer.config /etc/kafka/producer_newconfig.properties --whitelist 'datacollection_php_new' --num.streams 4
log_stdout=true
log_stderr=true
logfile=/var/log/kafka/kafkamirrorSrverDataCollection.out
logfile_maxbytes=20MB
logfile_backups=10
startsecs=10
startretries=2
priority=999
stopasgroup=true
user=netmining
environment=KAFKA_HEAP_OPTS="-Xmx1g -Xms1g",KAFKA_JMX_OPTS="-Dcom.sun.management.jmxremote -Dcom.sun.management.jmxremote.port=1093 -Dcom.sun.management.jmxremote.authenticate=false -Dcom.sun.management.jmxremote.ssl=false -Djava.rmi.server.hostname=<%= @ipaddress %>"
[program:kafka-production-wdc]
autorestart = true
autostart = true
directory = /usr/bin
command = kafka-mirror-maker --consumer.config /etc/kafka/consumer_newconfig.properties --producer.config /etc/kafka/producer_newconfig.properties --whitelist 'conversions_php|dynamic_creative_json|optout_json|wp_score_json|cookiematch_json|third_party_ids_json|wp_visitors_sessions_json' --num.streams 8
log_stdout=true
log_stderr=true
logfile=/var/log/kafka/kafkamirrorSrverProd.out
logfile_maxbytes=20MB
logfile_backups=10
startsecs=10
startretries=2
priority=999
stopasgroup=true
user=netmining
environment=KAFKA_HEAP_OPTS="-Xmx1g -Xms1g",KAFKA_JMX_OPTS="-Dcom.sun.management.jmxremote -Dcom.sun.management.jmxremote.port=1094 -Dcom.sun.management.jmxremote.authenticate=false -Dcom.sun.management.jmxremote.ssl=false -Djava.rmi.server.hostname=<%= @ipaddress %>"
[program:kafka-video_event_json-wdc]
autorestart = true
autostart = true
directory = /usr/bin
command = kafka-mirror-maker --consumer.config /etc/kafka/consumer_newconfig.properties --producer.config /etc/kafka/producer_newconfig.properties --whitelist 'video_event_json' --num.streams 2
log_stdout=true
log_stderr=true
logfile=/var/log/kafka/kafkamirrorVideoEvent.out
logfile_maxbytes=20MB
logfile_backups=10
startsecs=10
startretries=2
priority=999
stopasgroup=true
user=netmining
environment=KAFKA_HEAP_OPTS="-Xmx512M -Xms512M",KAFKA_JMX_OPTS="-Dcom.sun.management.jmxremote -Dcom.sun.management.jmxremote.port=1103 -Dcom.sun.management.jmxremote.authenticate=false -Dcom.sun.management.jmxremote.ssl=false -Djava.rmi.server.hostname=<%= @ipaddress %>"
<%- elsif @location == 'us-ash' -%>
[group:kafka-mirror-wdc]
programs=kafka-allbids-wdc,kafka-datacollection-wdc,kafka-production-wdc,kafka-video_event_json-wdc
priority=999
[program:kafka-allbids-wdc]
autorestart = true
autostart = true
directory = /usr/bin
command = kafka-mirror-maker --consumer.config /etc/kafka/consumer_newconfig.properties --producer.config /etc/kafka/producer_newconfig.properties --whitelist 'allbids_json' --num.streams 8
log_stdout=true
log_stderr=true
logfile=/var/log/kafka/kafkamirrorSrverAllBids.out
logfile_maxbytes=20MB
logfile_backups=10
startsecs=10
startretries=2
priority=999
stopasgroup=true
user=netmining
environment=KAFKA_HEAP_OPTS="-Xmx1g -Xms1g",KAFKA_JMX_OPTS="-Dcom.sun.management.jmxremote -Dcom.sun.management.jmxremote.port=1092 -Dcom.sun.management.jmxremote.authenticate=false -Dcom.sun.management.jmxremote.ssl=false -Djava.rmi.server.hostname=<%= @ipaddress %>"
[program:kafka-datacollection-wdc]
autorestart = true
autostart = true
directory = /usr/bin
command = kafka-mirror-maker --consumer.config /etc/kafka/consumer_newconfig.properties --producer.config /etc/kafka/producer_newconfig.properties --whitelist 'datacollection_php_new' --num.streams 4
log_stdout=true
log_stderr=true
logfile=/var/log/kafka/kafkamirrorSrverDataCollection.out
logfile_maxbytes=20MB
logfile_backups=10
startsecs=10
startretries=2
priority=999
stopasgroup=true
user=netmining
environment=KAFKA_HEAP_OPTS="-Xmx1g -Xms1g",KAFKA_JMX_OPTS="-Dcom.sun.management.jmxremote -Dcom.sun.management.jmxremote.port=1093 -Dcom.sun.management.jmxremote.authenticate=false -Dcom.sun.management.jmxremote.ssl=false -Djava.rmi.server.hostname=<%= @ipaddress %>"
[program:kafka-production-wdc]
autorestart = true
autostart = true
directory = /usr/bin
command = kafka-mirror-maker --consumer.config /etc/kafka/consumer_newconfig.properties --producer.config /etc/kafka/producer_newconfig.properties --whitelist 'conversions_php|dynamic_creative_json|optout_json|winnotification_json|wp_score_json|cookiematch_json|third_party_ids_json|wp_visitors_sessions_json' --num.streams 8
log_stdout=true
log_stderr=true
logfile=/var/log/kafka/kafkamirrorSrverProd.out
logfile_maxbytes=20MB
logfile_backups=10
startsecs=10
startretries=2
priority=999
stopasgroup=true
user=netmining
environment=KAFKA_HEAP_OPTS="-Xmx1g -Xms1g",KAFKA_JMX_OPTS="-Dcom.sun.management.jmxremote -Dcom.sun.management.jmxremote.port=1094 -Dcom.sun.management.jmxremote.authenticate=false -Dcom.sun.management.jmxremote.ssl=false -Djava.rmi.server.hostname=<%= @ipaddress %>"
[program:kafka-video_event_json-wdc]
autorestart = true
autostart = true
directory = /usr/bin
command = kafka-mirror-maker --consumer.config /etc/kafka/consumer_newconfig.properties --producer.config /etc/kafka/producer_newconfig.properties --whitelist 'video_event_json' --num.streams 2
log_stdout=true
log_stderr=true
logfile=/var/log/kafka/kafkamirrorVideoEvent.out
logfile_maxbytes=20MB
logfile_backups=10
startsecs=10
startretries=2
priority=999
stopasgroup=true
user=netmining
environment=KAFKA_HEAP_OPTS="-Xmx512M -Xms512M",KAFKA_JMX_OPTS="-Dcom.sun.management.jmxremote -Dcom.sun.management.jmxremote.port=1103 -Dcom.sun.management.jmxremote.authenticate=false -Dcom.sun.management.jmxremote.ssl=false -Djava.rmi.server.hostname=<%= @ipaddress %>"
<%- elsif @location == 'us-wdc' -%>
[group:kafka-mirror-wdc]
programs=kafka-allbids-wdc,kafka-winnotification-wdc,kafka-datacollection-wdc,kafka-production-wdc,kafka-profile_insights-wdc,kafka-video_event_json-wdc,kafka-creative_event_tracker_json-wdc
priority=999
[program:kafka-allbids-wdc]
autorestart = true
autostart = true
directory = /usr/bin
command = kafka-mirror-maker --consumer.config /etc/kafka/consumer_wdc_production.properties --producer.config /etc/kafka/producer_production.properties --whitelist 'allbids_json' --num.streams 4
log_stdout=true
log_stderr=true
logfile=/var/log/kafka/kafkamirrorSrverAllBids.out
logfile_maxbytes=20MB
logfile_backups=10
startsecs=10
startretries=2
priority=999
stopasgroup=true
user=netmining
environment=KAFKA_JMX_OPTS="-Dcom.sun.management.jmxremote.port=1092 -Dcom.sun.management.jmxremote.authenticate=false -Dcom.sun.management.jmxremote.ssl=false"
[program:kafka-datacollection-wdc]
autorestart = true
autostart = true
directory = /usr/bin
command = kafka-mirror-maker --consumer.config /etc/kafka/consumer_wdc_production.properties --producer.config /etc/kafka/producer_production.properties --whitelist 'datacollection_php_new' --num.streams 2
log_stdout=true
log_stderr=true
logfile=/var/log/kafka/kafkamirrorSrverDataCollection.out
logfile_maxbytes=20MB
logfile_backups=10
startsecs=10
startretries=2
priority=999
stopasgroup=true
user=netmining
environment=KAFKA_JMX_OPTS="-Dcom.sun.management.jmxremote.port=1093 -Dcom.sun.management.jmxremote.authenticate=false -Dcom.sun.management.jmxremote.ssl=false"
[program:kafka-production-wdc]
autorestart = true
autostart = true
directory = /usr/bin
command = kafka-mirror-maker --consumer.config /etc/kafka/consumer_wdc_production.properties --producer.config /etc/kafka/producer_production.properties --whitelist 'conversions_php|dynamic_creative_json|optout_json|wp_score_json|cookiematch_json|third_party_ids_json|wp_visitors_sessions_json|idlink_json' --num.streams 3
log_stdout=true
log_stderr=true
logfile=/var/log/kafka/kafkamirrorSrverProd.out
logfile_maxbytes=20MB
logfile_backups=10
startsecs=10
startretries=2
priority=999
stopasgroup=true
user=netmining
environment=KAFKA_JMX_OPTS="-Dcom.sun.management.jmxremote.port=1094 -Dcom.sun.management.jmxremote.authenticate=false -Dcom.sun.management.jmxremote.ssl=false"
[program:kafka-profile_insights-wdc]
autorestart = true
autostart = true
directory = /usr/bin
command = kafka-mirror-maker --consumer.config /etc/kafka/consumer_wdc_production.properties --producer.config /etc/kafka/producer_production.properties --whitelist 'pi_winnotification_json|pi_allbids_json|pi_reduced_bidrequests_json' --num.streams 4
log_stdout=true
log_stderr=true
logfile=/var/log/kafka/kafkamirrorProfileInsights.out
logfile_maxbytes=20MB
logfile_backups=10
startsecs=10
startretries=2
priority=999
stopasgroup=true
user=netmining
environment=KAFKA_JMX_OPTS="-Dcom.sun.management.jmxremote.port=1097 -Dcom.sun.management.jmxremote.authenticate=false -Dcom.sun.management.jmxremote.ssl=false"
[program:kafka-video_event_json-wdc]
autorestart = true
autostart = true
directory = /usr/bin
command = kafka-mirror-maker --consumer.config /etc/kafka/consumer_wdc_production.properties --producer.config /etc/kafka/producer_production.properties --whitelist 'video_event_json' --num.streams 2
log_stdout=true
log_stderr=true
logfile=/var/log/kafka/kafkamirrorVideoEvent.out
logfile_maxbytes=20MB
logfile_backups=10
startsecs=10
startretries=2
priority=999
stopasgroup=true
user=netmining
environment=KAFKA_JMX_OPTS="-Dcom.sun.management.jmxremote.port=1102 -Dcom.sun.management.jmxremote.authenticate=false -Dcom.sun.management.jmxremote.ssl=false"
[program:kafka-creative_event_tracker_json-wdc]
autorestart = true
autostart = true
directory = /usr/bin
command = kafka-mirror-maker --consumer.config /etc/kafka/consumer_wdc_production.properties --producer.config /etc/kafka/producer_production.properties --whitelist 'creative_event_tracker_json' --num.streams 2
log_stdout=true
log_stderr=true
logfile=/var/log/kafka/kafkamirrorCreativeEventTracker.out
logfile_maxbytes=20MB
logfile_backups=10
startsecs=10
startretries=2
priority=999
stopasgroup=true
user=netmining
environment=KAFKA_JMX_OPTS="-Dcom.sun.management.jmxremote.port=1103 -Dcom.sun.management.jmxremote.authenticate=false -Dcom.sun.management.jmxremote.ssl=false"
[program:kafka-winnotification-wdc]
autorestart = true
autostart = true
directory = /usr/bin
command = kafka-mirror-maker --consumer.config /etc/kafka/consumer_newconfig.properties --producer.config /etc/kafka/producer_wdc_agg_winnotification.properties --whitelist 'winnotification_json' --num.streams 4
log_stdout=true
log_stderr=true
logfile=/var/log/kafka/kafkamirrorWinnotification.out
logfile_maxbytes=20MB
logfile_backups=10
startsecs=10
startretries=2
priority=999
stopasgroup=true
user=netmining
environment=KAFKA_HEAP_OPTS="-Xmx1g -Xms1g",KAFKA_JMX_OPTS="-Dcom.sun.management.jmxremote.port=1104 -Dcom.sun.management.jmxremote.authenticate=false -Dcom.sun.management.jmxremote.ssl=false"
<%- else -%>
[group:kafka-mirror-wdc]
programs=kafka-allbids-wdc,kafka-datacollection-wdc,kafka-production-wdc,kafka-profile_insights-wdc,kafka-video_event_json-wdc,kafka-creative_event_tracker_json-wdc
priority=999
[program:kafka-allbids-wdc]
autorestart = true
autostart = true
directory = /usr/bin
command = kafka-mirror-maker --consumer.config /etc/kafka/consumer_wdc_production.properties --producer.config /etc/kafka/producer_production.properties --whitelist 'allbids_json' --num.streams 4
log_stdout=true
log_stderr=true
logfile=/var/log/kafka/kafkamirrorSrverAllBids.out
logfile_maxbytes=20MB
logfile_backups=10
startsecs=10
startretries=2
priority=999
stopasgroup=true
user=netmining
environment=KAFKA_JMX_OPTS="-Dcom.sun.management.jmxremote.port=1092 -Dcom.sun.management.jmxremote.authenticate=false -Dcom.sun.management.jmxremote.ssl=false"
[program:kafka-datacollection-wdc]
autorestart = true
autostart = true
directory = /usr/bin
command = kafka-mirror-maker --consumer.config /etc/kafka/consumer_wdc_production.properties --producer.config /etc/kafka/producer_production.properties --whitelist 'datacollection_php_new' --num.streams 2
log_stdout=true
log_stderr=true
logfile=/var/log/kafka/kafkamirrorSrverDataCollection.out
logfile_maxbytes=20MB
logfile_backups=10
startsecs=10
startretries=2
priority=999
stopasgroup=true
user=netmining
environment=KAFKA_JMX_OPTS="-Dcom.sun.management.jmxremote.port=1093 -Dcom.sun.management.jmxremote.authenticate=false -Dcom.sun.management.jmxremote.ssl=false"
[program:kafka-production-wdc]
autorestart = true
autostart = true
directory = /usr/bin
command = kafka-mirror-maker --consumer.config /etc/kafka/consumer_wdc_production.properties --producer.config /etc/kafka/producer_production.properties --whitelist 'conversions_php|dynamic_creative_json|optout_json|winnotification_json|wp_score_json|cookiematch_json|third_party_ids_json|wp_visitors_sessions_json|idlink_json' --num.streams 3
log_stdout=true
log_stderr=true
logfile=/var/log/kafka/kafkamirrorSrverProd.out
logfile_maxbytes=20MB
logfile_backups=10
startsecs=10
startretries=2
priority=999
stopasgroup=true
user=netmining
environment=KAFKA_JMX_OPTS="-Dcom.sun.management.jmxremote.port=1094 -Dcom.sun.management.jmxremote.authenticate=false -Dcom.sun.management.jmxremote.ssl=false"
[program:kafka-profile_insights-wdc]
autorestart = true
autostart = true
directory = /usr/bin
command = kafka-mirror-maker --consumer.config /etc/kafka/consumer_wdc_production.properties --producer.config /etc/kafka/producer_production.properties --whitelist 'pi_winnotification_json|pi_allbids_json|pi_reduced_bidrequests_json' --num.streams 4
log_stdout=true
log_stderr=true
logfile=/var/log/kafka/kafkamirrorProfileInsights.out
logfile_maxbytes=20MB
logfile_backups=10
startsecs=10
startretries=2
priority=999
stopasgroup=true
user=netmining
environment=KAFKA_JMX_OPTS="-Dcom.sun.management.jmxremote.port=1097 -Dcom.sun.management.jmxremote.authenticate=false -Dcom.sun.management.jmxremote.ssl=false"
[program:kafka-video_event_json-wdc]
autorestart = true
autostart = true
directory = /usr/bin
command = kafka-mirror-maker --consumer.config /etc/kafka/consumer_wdc_production.properties --producer.config /etc/kafka/producer_production.properties --whitelist 'video_event_json' --num.streams 2
log_stdout=true
log_stderr=true
logfile=/var/log/kafka/kafkamirrorVideoEvent.out
logfile_maxbytes=20MB
logfile_backups=10
startsecs=10
startretries=2
priority=999
stopasgroup=true
user=netmining
environment=KAFKA_JMX_OPTS="-Dcom.sun.management.jmxremote.port=1102 -Dcom.sun.management.jmxremote.authenticate=false -Dcom.sun.management.jmxremote.ssl=false"
[program:kafka-creative_event_tracker_json-wdc]
autorestart = true
autostart = true
directory = /usr/bin
command = kafka-mirror-maker --consumer.config /etc/kafka/consumer_wdc_production.properties --producer.config /etc/kafka/producer_production.properties --whitelist 'creative_event_tracker_json' --num.streams 2
log_stdout=true
log_stderr=true
logfile=/var/log/kafka/kafkamirrorCreativeEventTracker.out
logfile_maxbytes=20MB
logfile_backups=10
startsecs=10
startretries=2
priority=999
stopasgroup=true
user=netmining
environment=KAFKA_JMX_OPTS="-Dcom.sun.management.jmxremote.port=1103 -Dcom.sun.management.jmxremote.authenticate=false -Dcom.sun.management.jmxremote.ssl=false"
<%- end -%>
