# puppet/supervisor/conf.d/kafka-mirror-sao.conf

; /etc/supervisor/conf.d/kafka-supervisor-mirror--production.conf
;
; This file is managed by Puppet, do not edit locally !
<%- if @location == 'zg-ash' -%>
[group:kafka-mirror-sao]
programs=kafka-production-sao
priority=999
[program:kafka-production-sao]
autorestart = true
autostart = true
directory = /usr/bin
command = kafka-mirror-maker --consumer.config /etc/kafka/consumer_newconfig_zg_sao.properties --producer.config /etc/kafka/producer_newconfig.properties  --whitelist 'wp_score_json|wp_visitors_sessions_json|conversions_php' --num.streams 3
log_stdout=true
log_stderr=true
logfile=/var/log/kafka/kafkamirrorSrverProd.out
logfile_maxbytes=20MB
logfile_backups=10
startsecs=10
startretries=2
priority=999
stopasgroup=true
user=netmining
environment=KAFKA_JMX_OPTS="-Dcom.sun.management.jmxremote.port=1088 -Dcom.sun.management.jmxremote.authenticate=false -Dcom.sun.management.jmxremote.ssl=false"
<%- end -%>
<%- if @location == 'us-ash' -%>
[group:kafka-mirror-sao]
programs=kafka-production-sao
priority=999
[program:kafka-production-sao]
autorestart = true
autostart = true
directory = /usr/bin
command = kafka-mirror-maker --consumer.config /etc/kafka/consumer_sao_production.properties --producer.config /etc/kafka/producer_production.properties  --whitelist 'wp_score_json|wp_visitors_sessions_json|conversions_php' --num.streams 3
log_stdout=true
log_stderr=true
logfile=/var/log/kafka/kafkamirrorSrverProd.out
logfile_maxbytes=20MB
logfile_backups=10
startsecs=10
startretries=2
priority=999
stopasgroup=true
user=netmining
environment=KAFKA_JMX_OPTS="-Dcom.sun.management.jmxremote.port=1088 -Dcom.sun.management.jmxremote.authenticate=false -Dcom.sun.management.jmxremote.ssl=false"
<%- end -%>
<%- if @location == 'us-wdc' -%>
[group:kafka-mirror-sao]
programs=kafka-production-sao
priority=999
[program:kafka-production-sao]
autorestart = true
autostart = true
directory = /usr/bin
command = kafka-mirror-maker --consumer.config /etc/kafka/consumer_sao_production.properties --producer.config /etc/kafka/producer_production.properties  --whitelist 'wp_score_json|wp_visitors_sessions_json' --num.streams 3
log_stdout=true
log_stderr=true
logfile=/var/log/kafka/kafkamirrorSrverProd.out
logfile_maxbytes=20MB
logfile_backups=10
startsecs=10
startretries=2
priority=999
stopasgroup=true
user=netmining
environment=KAFKA_JMX_OPTS="-Dcom.sun.management.jmxremote.port=1088 -Dcom.sun.management.jmxremote.authenticate=false -Dcom.sun.management.jmxremote.ssl=false"
<%- end -%>
