# puppet/supervisor/conf.d/kafka-mirror-ams-fra.conf

; /etc/supervisor/conf.d/kafka-supervisor-mirror--production.conf
;
; This file is managed by Puppet, do not edit locally !
[group:kafka-mirror-ams-fra]
programs=kafka-datacollection-ams-fra,kafka-production-ams-fra
priority=999
[program:kafka-datacollection-ams-fra]
autorestart = true
autostart = true
directory = /usr/bin
command = kafka-mirror-maker --consumer.config /etc/kafka/consumer_ams_fra_production.properties --producer.config /etc/kafka/producer_production.properties  --whitelist 'datacollection_php_new' --num.streams 2
log_stdout=true
log_stderr=true
logfile=/var/log/kafka/kafkamirrorSrverDataCollection.out
logfile_maxbytes=20MB
logfile_backups=10
startsecs=10
startretries=2
priority=999
stopasgroup=true
user=netmining
environment=KAFKA_JMX_OPTS="-Dcom.sun.management.jmxremote.port=1082 -Dcom.sun.management.jmxremote.authenticate=false -Dcom.sun.management.jmxremote.ssl=false"
[program:kafka-production-ams-fra]
autorestart = true
autostart = true
directory = /usr/bin
command = kafka-mirror-maker --consumer.config /etc/kafka/consumer_ams_fra_production.properties --producer.config /etc/kafka/producer_production.properties  --whitelist 'conversions_php|dynamic_creative_json|wp_score_json|cookiematch_json|wp_visitors_sessions_json' --num.streams 3
log_stdout=true
log_stderr=true
logfile=/var/log/kafka/kafkamirrorSrverProd.out
logfile_maxbytes=20MB
logfile_backups=10
startsecs=10
startretries=2
priority=999
stopasgroup=true
user=netmining
environment=KAFKA_JMX_OPTS="-Dcom.sun.management.jmxremote.port=1083 -Dcom.sun.management.jmxremote.authenticate=false -Dcom.sun.management.jmxremote.ssl=false"
