# puppet/supervisor/conf.d/mmmhome.conf
; /etc/supervisor/conf.d/mmmhome.conf
;
; This file is managed by Puppet, do not edit locally !

[program:mmmhome]
autorestart = true
autostart = true
command = /usr/bin/python MMM_OOR_API.py
directory = /var/ignitionone/projects/rserver/MMMhome
redirect_stderr = true
stdout_logfile=/var/log/ignitionone/mmmhome.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=10
startsecs=10
startretries=999
user=ignitionone
