# puppet/supervisor/conf.d/process-exporter-worker.conf
; /etc/supervisor/conf.d/process-exporter-worker.conf
;
; This file is managed by Puppet, do not edit locally !
<%- if ( @fqdn != 'vadmin-onsite-001.us-atl.ignitionone.com' ) -%>
<%- case @area
when "production", "staging"
-%>
[program:process-exporter-worker]
autorestart = true
autostart = true
directory = /home/netmining
numprocs = 5
process_name = process_exporter_worker_%(process_num)s
environment = SUPERVISOR_PROCESS='worker%(process_num)s'
command = /usr/bin/php /var/netmining/projects/Onsite/Exporter/console export:lead:handle process
user=netmining
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=1
<%- else -%>
[program:process-exporter-preprod-worker]
autorestart = true
autostart = false
directory = /home/netmining
numprocs = 3
process_name = process_exporter_preprod_worker_%(process_num)s
environment = SUPERVISOR_PROCESS='worker%(process_num)s'
command = /usr/bin/php /var/netmining/projects/Onsite/Exporter/console export:lead:handle process
user=netmining
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=1
<%- end -%>
<%- end -%>


