# puppet/supervisor/conf.d/secor.conf
; /etc/supervisor/conf.d/secor.conf
;
; This file is managed by Puppet, do not edit locally !

<%- if (!@secor_config.nil?) then -%>
<%- @secor_config.each do |config_name,config_hash| -%>
[program:secor-<%= config_name %>]
autorestart = true
autostart = true
directory = /opt/secor/bin
command = java -ea -Dsecor_group=secor_backup -Dlog4j.configuration=file:///etc/secor/config/log4j.prod.properties -Dconfig=file:///etc/secor/config/secor.prod.<%= config_name %>.properties -cp secor-0.16-SNAPSHOT.jar:lib/* com.pinterest.secor.main.ConsumerMain
log_stdout=false
log_stderr=true
logfile=/var/log/secor/secor.out
logfile_maxbytes=20MB
logfile_backups=10
startsecs=10
startretries=2
priority=999
stopasgroup=true
<%= %>
<%- end -%>
<%- end -%>
