# puppet/supervisor/conf.d/storm-ui-uat.conf
; /etc/supervisor/conf.d/storm-ui.conf
;
; This file is managed by Puppet, do not edit locally !

[program:storm-ui]
autorestart = true
autostart = true
directory = /opt/storm/bin
command = /opt/storm/bin/storm -c /etc/default/storm-ui ui
log_stdout=true
log_stderr=true
logfile=/var/log/storm/ui.out
logfile_maxbytes=20MB
logfile_backups=10
startsecs=10
startretries=2
user=storm
environment=STORM_CONF_DIR="/etc/storm/"
