# puppet/supervisor/conf.d/storm-nimbus-uat.conf
; /etc/supervisor/conf.d/storm-nimbus.conf
;
; This file is managed by Puppet, do not edit locally !

[program:storm-nimbus]
autorestart = true
autostart = true
directory = /opt/storm/bin
command = /opt/storm/bin/storm -c /etc/default/storm-nimbus nimbus
log_stdout=true
log_stderr=true
logfile=/var/log/storm/nimbus.out
logfile_maxbytes=20MB
logfile_backups=10
startsecs=10
startretries=2
user=storm
environment=STORM_CONF_DIR="/etc/storm/"
